# com.vives.bird-pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vives.bird-pipeline</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/vives/bird-classifier/venv-coral/bin/python3</string>
        <string>-u</string>
        <string>/Users/vives/bird-classifier/bird_pipeline.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/vives/bird-classifier</string>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/vives/bird-snapshots/logs/pipeline.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/vives/bird-snapshots/logs/pipeline-stderr.log</string>
</dict>
</plist>
