# deploy/systemd/bird-dashboard.service
[Unit]
Description=Bird dashboard — FastAPI via uvicorn
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/bird-classifier
Environment=PI_MODE=1
ExecStart=%h/bird-classifier/venv/bin/uvicorn dashboard.api:app --host 0.0.0.0 --port 8099
Restart=always
RestartSec=5
StandardOutput=append:%h/logs/bird-dashboard.log
StandardError=append:%h/logs/bird-dashboard.log

[Install]
WantedBy=default.target
