# deploy/systemd/service-canary.timer
[Unit]
Description=Run service-canary every 2 min
[Timer]
OnBootSec=3min
OnUnitActiveSec=2min
[Install]
WantedBy=timers.target
