# deploy/systemd/service-canary.service
[Unit]
Description=Zombie self-heal canary (dashboard + sshd)
[Service]
Type=oneshot
ExecStart=/usr/local/sbin/service-canary.sh
